# docs/ocr-service.service
[Unit]
Description=OCR Kataloogijälgimise Teenus (Qwen VLM)
After=network.target

[Service]
Type=simple
User=mf
WorkingDirectory=/home/mf/Dokumendid/LLM/qwen-treening
Environment="PATH=/home/mf/Dokumendid/LLM/qwen-treening/venv/bin:/usr/local/bin:/usr/bin"
ExecStart=/home/mf/Dokumendid/LLM/qwen-treening/venv/bin/python /home/mf/Dokumendid/LLM/qwen-treening/kataloogi-jalgimine-ja-ocr.py
Restart=on-failure
RestartSec=10

# GPU ja mälu
Nice=-5
MemoryMax=90%

# Logid (lisaks skripti oma logifailile)
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
